# files/etc/systemd/system/kubelet.service.d/10-extra-args.conf
[Service]
Environment="KUBELET_EXTRA_ARGS=--container-runtime=remote --container-runtime-endpoint=unix:///var/run/crio/crio.sock --cgroup-driver=systemd --config=/etc/kubernetes/kubelet-config.yaml"
